# templates/program.conf.j2
[program:{{ item.name }}]
command={{ item.command }}
